# system_files/common/usr/lib/systemd/system/ostree-finalize-staged.service.d/50-cache22-uki.conf
# cache22 drop-in: build + sign UKIs after ostree finalizes a staged
# deploy. systemd runs ExecStop= lines in declaration order, so this
# fires AFTER `ostree admin finalize-staged` writes the BLS entry —
# which resign-uki needs (it reads cmdline + kernel paths from the
# BLS entry written by finalize).
#
# Same job as ostree-finalize-staged, so:
#   - same shutdown-blocking semantics (system waits for UKI build
#     before reboot/poweroff/halt/soft-reboot)
#   - failure recorded by ostree's existing machinery
#   - covers every shutdown verb (reboot/poweroff/halt/kexec/soft-reboot)
#     because the parent unit's Conflicts=final.target catches all of them

[Unit]
RequiresMountsFor=/efi

[Service]
ExecStop=/usr/libexec/cache22/resign-uki
TimeoutStopSec=10m
