# system_files/common/usr/lib/systemd/system/ostree-remount.service.d/50-cache22-etc-rw.conf
# cache22 drop-in: ensure /etc is writable post-soft-reboot.
#
# ostree-remount.service handles /var and /sysroot remounts and a few
# transient-etc cases (when prepare-root.conf has [composefs] enabled
# and etc.transient=true). It does NOT re-establish /etc as a
# writable bind for the legacy non-composefs / non-transient case.
#
# Why this matters: systemd's soft-reboot pivot only preserves a
# fixed list of mounts (/var, /boot, /efi, /sysroot, /home, etc.).
# /etc isn't on that list, so our prepare-soft-reboot's bind on
# /etc gets dropped during the pivot, leaving /etc as a plain dir
# on the (RO btrfs) deploy mount → /etc is RO in the new userspace.
#
# By piggybacking on ostree-remount via ExecStartPost, we inherit the
# right ordering for free: After=local-fs (so / is mounted) and
# Before=systemd-tmpfiles-setup.service / systemd-rfkill.service /
# local-fs.target (so anything that needs writable /etc waits for us).
# Same systemd job, idempotent on hard boot (where ostree-prepare-root
# in initrd already established the bind).

[Service]
ExecStartPost=/usr/libexec/cache22/ensure-etc-writable
