# system_files/common/usr/lib/systemd/system/cache22-healthcheck.timer
[Unit]
Description=Run cache22 health check 2 minutes after boot

[Timer]
OnBootSec=2min
Unit=cache22-healthcheck.service
RemainAfterElapse=true

[Install]
WantedBy=timers.target
