# system_files/common/usr/lib/systemd/system/cache22-healthcheck.service
[Unit]
Description=cache22 post-boot health check (auto-rollback after repeated failures)
ConditionPathExists=/etc/cache22/healthcheck.d/required.d
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/cache22-healthcheck
StandardOutput=journal
StandardError=journal+console
